# dasmon/server/dasmond.conf
# THIS FILE CREATED BY CFENGINE
# DO NOT EDIT DIRECTLY

# Dasmon Daemon - Upstart Init Script
#    - installed in /etc/init/dasmond.conf

# Upstart Commands:
#    - to start:     /sbin/start dasmond
#    - to stop:      /sbin/stop dasmond
#    - to restart:   /sbin/restart dasmond
#    - for status:   /sbin/status dasmond

start on stopped rc or RUNLEVEL=[2345]
stop on runlevel [!2345]

#respawn limit 2 10
respawn
exec su - dasmon -s /bin/bash -c "/usr/local/bin/dasmond"

