# dasmon/server/init.dasmond.conf
# THIS FILE CREATED BY CFENGINE
# DO NOT EDIT DIRECTLY

start on stopped rc or RUNLEVEL=[2345]
stop on runlevel [!2345]

#respawn limit 2 10
respawn
exec su - dasmon -s /bin/bash -c "/usr/local/bin/dasmond"
