# dasmon/client/dasmongui.conf
# THIS FILE CREATED BY CFENGINE
# DO NOT EDIT DIRECTLY

# Dasmon Client GUI - Upstart Init Script
#    - installed in /etc/init/dasmongui.conf

# Upstart Commands:
#    - to start:     /sbin/start dasmongui
#    - to stop:      /sbin/stop dasmongui
#    - to restart:   /sbin/restart dasmongui
#    - for status:   /sbin/status dasmongui

# Starts the Dasmon GUI on a dedicated X server
#    - calls XInit on Dasmon Client GUI configuration in:
#       -> /etc/X11/xinit/xinitrc.dasmon

start on stopped rc or RUNLEVEL=4
stop on runlevel [!4]

respawn limit 2 10
exec /usr/bin/xinit /etc/X11/xinit/xinitrc.dasmon

