# dasmon/client/init.dasmongui.conf
# THIS FILE CREATED BY CFENGINE
# DO NOT EDIT DIRECTLY

#Starts the dasmon GUI on a dedicated X server

start on stopped rc or RUNLEVEL=4
stop on runlevel [!4]

respawn limit 2 10
exec /usr/bin/xinit /etc/X11/xinit/xinitrc.dasmon
